# hotbox/app_docker_template/inittab
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

ttyS0::respawn:/bin/ash /start.sh

::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/openrc shutdown
